# out/modbus_module.service
[Unit]
 Description=WolkAbout Gateway Module for Modbus devices
 After=multi-user.target

[Service]
 Type=idle
 WorkingDirectory=/usr/local/bin/
 ExecStart=/usr/local/bin/modbusModule /etc/modbusModule/moduleConfiguration.json /etc/modbusModule/devicesConfiguration.json
 Restart=always

[Install]
 WantedBy=multi-user.target
